# etc/supervisord.d/ipfs.ini
[program:ipfs]
directory=/data/
user=root
command=/opt/qnib/ipfs/bin/start.sh
autostart=true
exitcodes=0
startsecs=3
stopwaitsecs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
